# sensor/plantmon-sensor.service
[Unit]
Description=plantmon sensor service
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/plantmon-sensor
